# macOs/macsetup/com.user.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOUR_USERNAME/backup_files.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>86400</integer> <!-- Runs every 24 hours -->

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/backup_log.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/backup_log.err</string>
</dict>
</plist>

